# kafka/files/kafka-mirror.service
# Generated by Saltstack --  All edits will be lost !!

[Unit]
Description=Apache Kafka (Mirror Maker)

[Service]
User={{ user }}
Group={{ group }}
EnvironmentFile=/etc/default/{{ service }}
ExecStart=/bin/bash $KAFKA_RUNNER kafka.tools.MirrorMaker $KAFKA_ARGS $KAFKA_CONFIG
TimeoutSec=10
TimeoutStopSec={{ timeout }}
LimitNOFILE={{ ulimit }}
StandardOutput=/var/log/kafka/{{ service }}.log
StandardError=/var/log/kafka/{{ service }}_error.log

[Install]
WantedBy=multi-user.target
